# systemd/etc/systemd/system/docker-rbac-composed.service
[Unit]
Description=RBAC container
Requires=docker.service
After=docker.service

[Service]
Restart=always
#ExecStartPre=/bin/mkdir -p /usr/local/etc/docker-rbac-composed/
ExecStart=/usr/local/bin/docker-compose -f /usr/local/etc/docker-rbac-composed/docker-compose.yml up
ExecStop=/usr/local/bin/docker-compose -f /usr/local/etc/docker-rbac-composed/docker-compose.yml down

[Install]
WantedBy=multi-user.target
